# modules/supervisorprog/templates/program.conf.erb
[program:<%= @title %>]
process_name=%(program_name)s%(process_num)s
command=<%= @command %>
directory=<%= @directory %>
user=<%= @user %>
numprocs=<%= @numprocs %>
autostart=<%= @autostart %>
autorestart=true
redirect_stderr=true
environment=<%= @enviro %>
stdout_logfile=%(program_name)s_stdout.log
stderr_logfile=%(program_name)s_stderr.log
